# cepheus-lightbroker/deb/postinst
#!/bin/sh
update-rc.d -f cepheus-broker defaults
service cepheus-broker start
