# cepheus-cep/deb/postinst
#!/bin/sh
update-rc.d -f cepheus-cep defaults
service cepheus-cep start
