# cepheus-cep/deb/prerm
#!/bin/sh

service cepheus-cep stop

if [ "$1" = "remove" ]; then
    update-rc.d -f cepheus-cep remove
fi
